# build/analyzer/build/analyzer.conf
[program:analyzer]
command=/bin/sh -c 'if [ -f /app/config/analyzer.yaml ]; then /usr/local/bin/analyzer --config /app/config/analyzer.yaml; else /usr/local/bin/analyzer --listen-addr ${LISTEN_ADDR:-0.0.0.0:50052} --database-uri ${DATABASE_URI}; fi'
directory=/app
autostart=true
autorestart=true
startretries=3
user=%(ENV_BUILD_USER)s
redirect_stderr=true
stdout_logfile=/var/log/supervisor/analyzer.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=CONFIG_DIR="/app/config"
