# build/checker/build/checker.conf
[program:checker]
command=/bin/bash -lc "echo 'checker container started. run /app/scripts/quick_check.sh or /app/scripts/comprehensive_test.sh'; tail -f /dev/null"
directory=/app
autostart=true
autorestart=true
startretries=3
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/checker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
