# agent-build/build/agent.conf
[program:agent]
command=/bin/sh -c 'AGENT_DIR="${RPINGMESH_AGENT_DIR_PATH:-/app}"; CONFIG_FILE="${AGENT_DIR}/config/agent.yaml"; echo "=== Agent Command Debug Info ===" && echo "Timestamp: $(date -Iseconds)" && echo "RPINGMESH_AGENT_DIR_PATH=${AGENT_DIR}" && echo "Config File: ${CONFIG_FILE}" && echo "Controller Addr (env): ${RPINGMESH_CONTROLLER_ADDR:-<not set, will use config>}" && echo "Analyzer Addr (env): ${RPINGMESH_ANALYZER_ADDR:-<not set, will use config>}" && echo "Analyzer Enabled (env): ${RPINGMESH_ANALYZER_ENABLED:-<not set, default=false>}" && echo "Effective command:" && echo "  /usr/local/bin/agent --config ${CONFIG_FILE} ${RPINGMESH_CONTROLLER_ADDR:+--controller-addr ${RPINGMESH_CONTROLLER_ADDR}} ${RPINGMESH_ANALYZER_ADDR:+--analyzer-addr ${RPINGMESH_ANALYZER_ADDR}}" && echo "=== Starting Agent Process ===" && /usr/local/bin/agent --config "${CONFIG_FILE}" ${RPINGMESH_CONTROLLER_ADDR:+--controller-addr ${RPINGMESH_CONTROLLER_ADDR}} ${RPINGMESH_ANALYZER_ADDR:+--analyzer-addr ${RPINGMESH_ANALYZER_ADDR}}'
directory=/tmp
autostart=true
autorestart=true
startretries=3
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/agent.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5