# agent-build/build/agent.conf
[program:agent]
command=/usr/local/bin/agent --config /app/config/agent.yaml
# All configuration is loaded from agent.yaml file
# Environment variables (RPINGMESH_*) can optionally override config file values (Viper handles this automatically)
directory=/tmp
autostart=true
autorestart=true
startretries=3
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/agent.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5