# agent-build/build/agent.conf
[program:agent]
command=/bin/sh -c '/usr/local/bin/agent --config /app/config/agent.yaml ${RPINGMESH_CONTROLLER_ADDR:+--controller-addr ${RPINGMESH_CONTROLLER_ADDR}} ${RPINGMESH_ANALYZER_ADDR:+--analyzer-addr ${RPINGMESH_ANALYZER_ADDR}}'
directory=/app
autostart=true
autorestart=true
startretries=3
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/agent.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=CONFIG_DIR="/app/config"


