# controller-build/build/controller.conf
[program:controller]
command=/usr/local/bin/controller --config /app/config/controller.yaml
directory=/app
autostart=true
autorestart=true
startretries=3
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/controller.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=CONFIG_DIR="/app/config"
