# build/rqlite/build/rqlite.conf
[program:rqlite]
command=/bin/sh -c '/usr/local/bin/rqlited -node-id ${NODE_ID} -http-addr 0.0.0.0:4001 -http-adv-addr ${NODE_ID}:4001 -raft-addr 0.0.0.0:4002 -raft-adv-addr ${NODE_ID}:4002 ${JOIN_ADDR:+-join ${JOIN_ADDR}} /data'
directory=/data
autostart=true
autorestart=true
startretries=3
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/rqlite.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
